# management-ui/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/var/log/nginx/stdout.log
stderr_logfile=/var/log/nginx/stderr.log
autorestart=true
priority=10

[program:backend]
command=/usr/local/bin/uvicorn app:app --host 127.0.0.1 --port 8000 --workers 4
directory=/app/backend
stdout_logfile=/var/log/exapg/backend.log
stderr_logfile=/var/log/exapg/backend.log
autorestart=true
priority=20
environment=PYTHONPATH=/app/backend

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface 